# deploy/systemd/suppoint-bot-dev.service
[Unit]
Description=SUPPOINT Telegram Bot (development, autoreload)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=andrii
WorkingDirectory=/home/andrii/suppoint-bot
EnvironmentFile=/home/andrii/suppoint-bot/.env
ExecStart=/home/andrii/suppoint-bot/.venv/bin/watchmedo auto-restart --patterns=*.py --recursive -- /home/andrii/suppoint-bot/.venv/bin/python /home/andrii/suppoint-bot/bot.py
Restart=always
RestartSec=3
TimeoutStopSec=20
KillSignal=SIGINT
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
